# PKGBUILD/Predator-Sense-systemd-git.install
post_install() {

    sudo systemctl daemon-reload
    sudo systemctl start turbo-fan
    sudo systemctl enable turbo-fan

}

post_remove(){

    sudo systemctl stop turbo-fan
    sudo systemctl disable turbo-fan
    sudo systemctl daemon-reload
    sudo rm -rvf /opt/turbo-fan

}
